# server/rootfs/etc/periodic/15min/dso_tasks.sh
#!/bin/sh
set -e

s6-setuidgid nginx php84 /var/www/_tasks.php
