# server/rootfs/etc/s6/anubis/run
#!/bin/sh
set -e

mkdir -p /run/anubis
chown anubis:anubis /run/anubis

exec s6-setuidgid anubis /usr/bin/anubis -bind-network unix -bind /run/anubis/anubis.sock -target /run/caddy/internal.sock
