# server/rootfs/etc/s6/caddy/run
#!/bin/sh
set -e

export HOME=/var/lib/caddy

mkdir -p /run/caddy
chown caddy:caddy /run/caddy

exec s6-setuidgid caddy /usr/sbin/caddy run --config /etc/caddy/Caddyfile
