# server/rootfs/etc/s6/caddy/run
#!/bin/sh

export HOME=/var/lib/caddy

exec s6-setuidgid caddy /usr/sbin/caddy run --config /etc/Caddyfile
