# server/rootfs/etc/s6/php-fpm/run
#!/bin/sh
set -e

/bin/configure.sh root:php

exec /usr/sbin/php-fpm84 --nodaemonize --fpm-config /etc/php-fpm.conf
