# server/rootfs/etc/s6/php-fpm/run
#!/bin/sh

# BEGIN Create config.php
PRIVCONF=/run/dso-config.php
touch "${PRIVCONF}"
chown root:php "${PRIVCONF}"
chmod 640 "${PRIVCONF}"

phpenvvar() {
    VAR="${1}"
    eval VAL="\$$VAR"
    echo "\$${VAR} = '${VAL}';" >> "${PRIVCONF}"
}

phpenvvar_b64() {
    VAR="${1}"
    eval VAL="\$$VAR"
    echo "\$${VAR} = base64_decode('${VAL}');" >> "${PRIVCONF}"
}

echo '<?php' > "${PRIVCONF}"

phpenvvar 'DB_HOST'
phpenvvar 'DB_USERNAME'
phpenvvar 'DB_PASSWORD'
phpenvvar 'DB_DATABASE'

phpenvvar_b64 'JWT_PRIVATE_KEY'
phpenvvar_b64 'JWT_PUBLIC_KEY'

phpenvvar 'CAPTCHA_SECRET_KEY'
phpenvvar 'CAPTCHA_FONT'

echo "<?php require_once('${PRIVCONF}');" > /var/www/api/config.php
# END Create config.php

# BEGIN Create msmtp config
MSMTP_CONF=/run/msmtp.conf
touch "${MSMTP_CONF}"
chown root:php "${MSMTP_CONF}"
chmod 640 "${MSMTP_CONF}"

echo 'defaults' > "${MSMTP_CONF}"
echo 'auth on' >> "${MSMTP_CONF}"
echo 'tls on' >> "${MSMTP_CONF}"
echo 'tls_starttls on' >> "${MSMTP_CONF}"
echo "port 587" >> "${MSMTP_CONF}"
echo 'account default' >> "${MSMTP_CONF}"
echo "host ${SMTP_HOST}" >> "${MSMTP_CONF}"
echo "from ${SMTP_FROM}" >> "${MSMTP_CONF}"
echo "user ${SMTP_USERNAME}" >> "${MSMTP_CONF}"
echo "passwordeval echo '${SMTP_PASSWORD}'" >> "${MSMTP_CONF}"
# END Create msmtp config

exec /usr/sbin/php-fpm84 --nodaemonize --fpm-config /etc/php-fpm.conf
