# server/rootfs/etc/s6/php-fpm/run
#!/bin/sh
set -e

/bin/configure_msmtp.sh

cp -f /var/www/releases_fallback.json /tmp/releases.json
chown php:php /tmp/releases.json

exec /usr/sbin/php-fpm84 --nodaemonize --fpm-config /etc/php-fpm.conf
