# server/rootfs/etc/s6/php-fpm/run
#!/bin/sh
set -e

mkdir -p /run/darksignsonline
chmod 750 /run/darksignsonline
chown root:php /run/darksignsonline
/bin/configure.sh

cp -f /var/www/releases_fallback.json /run/darksignsonline/releases.json
chown php:php /run/darksignsonline/releases.json

exec /usr/sbin/php-fpm84 --nodaemonize --fpm-config /etc/php-fpm.conf
