# server/rootfs/etc/s6/crond/run
#!/bin/sh
exec /usr/sbin/crond -f -c /etc/crontabs -d 8
